# data/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

[program:watcherd]
command=watcherd -v -p /shared/httpd -a "create-vhost.sh '%%p' '%%n' '__TLD__' '%%p/__MASS_VHOST_TPL__/' '__CA_KEY__' '__CA_CRT__' '__VERBOSE__' '1'" -d "rm /etc/httpd/vhost.d/%%n.conf" -t "nginx -s reload"
#startsecs = 0
#autorestart = false
#startretries = 1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
